# audio_cache/supervisor.conf
[program:ai_server]
command=/home/qc/ai_server/ai_server/web -f conf/web.json
 # 启动命令
directory=/home/qc/ai_server/ai_server  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
autorestart=true  # 自动重启
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync1.5/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/go/go/bin"
stdout_logfile=/home/qc/log/ai_server.log  # 标准输出日志文件

[program:latent_master]
command=python -u -m server.master
 # 启动命令
directory=/home/qc/workspace/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
numprocs=1
autorestart=true  # 自动重启
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/go/go/bin"
stdout_logfile=/home/qc/log/latent.master.log  # 标准输出日志文件

[program:latent_affine]
command=python -u -m server.affine_transform --gpu_id=%(process_num)d
 # 启动命令
directory=/home/qc/workspace/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
numprocs=4
autorestart=true  # 自动重启
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/go/go/bin"
stdout_logfile=/home/qc/log/latent.affine.log  # 标准输出日志文件

[program:latent_affine]
command=python -u -m server.affine_train
 # 启动命令
directory=/home/qc/workspace/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
autorestart=true  # 自动重启
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/go/go/bin"
stdout_logfile=/home/qc/log/latent.affine_train.log  # 标准输出日志文件

[program:latent_inference]
command=python -u -m server.inference --gpu_id=%(process_num)d
 # 启动命令
directory=/home/qc/workspace/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
numprocs=4
autorestart=true  # 自动重启
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/go/go/bin"
stdout_logfile=/home/qc/log/latent.inference.log  # 标准输出日志文件

[program:latent_restore]
command=python -u -m server.restore --gpu_id=%(process_num)d
 # 启动命令
directory=/home/qc/workspace/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
numprocs=8
autorestart=true  # 自动重启
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/go/go/bin"
stdout_logfile=/home/qc/log/latent.restore.log  # 标准输出日志文件

[program:api_pipeline_fy]
command=python -u api_pipeline_fy.py
 # 启动命令
directory=/home/qc/train/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
autorestart=true  # 自动重启
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync1.5,PATH="/home/qc/miniconda3/envs/latentsync1.5/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/miniconda3/condabin"
stdout_logfile=/home/qc/log/fy_queue_pipeline.log  # 标准输出日志文件

[program:api_pipeline_audio]
command=python -u api_audio_pipeline.py
 # 启动命令
directory=/home/qc/workspace/LatentSync  # 工作目录
user=qc  # 运行程序的用户
autostart=true  # 自动启动
autorestart=true  # 自动重启
redirect_stderr=true
environment=CONDA_DEFAULT_ENV=latentsync,PATH="/home/qc/miniconda3/envs/latentsync/bin:/home/qc/miniconda3/condabin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games:/snap/bin:/home/qc/miniconda3/condabin:/home/qc/miniconda3/condabin"
stdout_logfile=/home/qc/log/audio_queue_pipeline.log  # 标准输出日志文件